# supervisord.conf
[supervisord]
nodaemon=true

[program:scheduler]
command=/bin/bash -c "rake resque:scheduler"

[program:puma]
command=/bin/bash -c "bundle exec puma -C ./config/puma.rb"
